# plex/rootfs/etc/services.d/plex/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Plex Media Server
# Take down the S6 supervision tree when the Plex Media Server fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "Plex crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "Plex stopped, restarting..."
